# scripts/openclaw-config-restart.service
[Unit]
Description=Restart OpenClaw after config change
Documentation=https://docs.openclaw.ai/

[Service]
Type=oneshot
ExecStart=/bin/systemctl restart openclaw.service
